# templates/default/systemd-geminabox.service.erb
[Unit]
Description=Geminabox Application
Documentation=http://docs.geminabox.com
After=network.target
Requires=

[Service]
Type=notify
<% if @tmpdir -%>
Environment="TMPDIR=<%= @tmpdir %>"
<% end -%>
ExecStart=<%= Docker::Helpers.executable(node) %> -d<%= @daemon_options %>

[Install]
Also=geminabox.socket